# runtimes/012_nginx_phpfpm_single/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autorestart=true
